# etc/crontab
*/1 * * * * root /src/script/reaptrees.sh /src
