# cron/crontab
*/1 * * * * root /src/reaptrees.sh /trees
